# piplug.service
[Unit]
Description=piplug

[Service]
Type=simple
ExecStart=/usr/bin/python /home/pi/piplug.py
Restart=on-abort

[Install]
WantedBy=multi-user.target
